# docker/rootfs/etc/cron.d/01_mangadlp
# default crontab to run manga-dlp once a day
# and get all (new) chapters of the mangas in
# the file mangas.txt

0 3 * * *   abc    python3 /app/manga-dlp.py --read /app/mangas.txt -c all --wait 2
